# rootfs/etc/cron.d/fwtree
*/5 *   * * * root ipset-helper
07  */4 * * * root fwtree-update-blacklists
