# rootfs/etc/cron.d/fwtree
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

*/5 *   * * * root ipset-helper
#07  */4 * * * root fwtree-update-blacklists
